# project/classics-backend.service
[Unit]
Description=Classics Learning Platform Backend
After=network.target mysql.service

[Service]
Type=exec
User=www-data
WorkingDirectory=/var/www/classics-learning/backend
Environment=PATH=/var/www/classics-learning/backend/venv/bin
EnvironmentFile=/var/www/classics-learning/backend/.env
ExecStart=/var/www/classics-learning/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
